# rc.local
#/bin/bash
bash /opt/damn/Damn-Vulnerable-Drone/start.lite.with.root.sh &


