# rc.local
#!/bin/bash
#ln -s /opt/Damn-Vulnerable-Drone/rc.local /etc/rc.local
sleep 30
bash /opt/Damn-Vulnerable-Drone/start.lite.with.root.sh >/dev/null 2>&1 &


exit 0
